# deploy/com.dog-finder.daily-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dog-finder.daily-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/fanguard/Code/dog-finder/scripts/daily-refresh.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/fanguard/Code/dog-finder/logs/daily-refresh.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fanguard/Code/dog-finder/logs/daily-refresh.launchd.log</string>
</dict>
</plist>
